# sys-fs/swapspace/files/swapspace.initd
#!/sbin/openrc-run

description="dynamic swap manager"
command="/usr/sbin/swapspace"
pidfile="/run/${RC_SVCNAME}.pid"
command_args="-p ${pidfile} --quiet"
command_args_background="--daemon"

depend() {
    need localmount
}
